# templates/sagemcom-3864-exporter.service.j2
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/python /opt/sagemcom_3864_exporter.py
User=prometheus
WorkingDirectory=/opt

